# etc/cron.d/redeploy
* * * * * root /usr/local/bin/redeploy-engine.sh >> /var/log/redeploy.log 2>&1

